# deploy/podman/roamarr.container
[Unit]
Description=Roamarr
Wants=network-online.target
After=network-online.target

[Container]
Image=localhost/roamarr:latest
Pull=never
ContainerName=roamarr
PublishPort=3000:3000
Volume=roamarr-data:/data
EnvironmentFile=%h/.config/containers/systemd/roamarr.env
Environment=NODE_ENV=production
Environment=DATABASE_PATH=/data/roamarr.db
Environment=PORT=3000

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
